# src/main/systemd/energieid.timer
[Unit]
Description=Daily EnergieID upload.

[Timer]
OnCalendar=*-*-* 21:55:00
Persistent=true

[Install]
WantedBy=timers.target
